# systemd/konsum.socket
[Socket]
ListenStream=/run/{{name}}/socket
MaxConnections=4
