# systemd/konsum.socket
[Socket]
ListenStream=/run/konsum/socket
MaxConnections=4
