# systemd/konsum.socket
[Socket]
ListenStream={{socket}}
MaxConnections=4
NoDelay=true
FileDescriptorName=api

[Install]
RequiredBy={{name}}.service
