# systemd/konsum.socket
[Socket]
ListenStream=/run/{{name}}/socket
MaxConnections=4
NoDelay=true

[Install]
RequiredBy={{name}}.service
